# directory-app/housekeeping-agent.service
[Unit]
Description=Polka Housekeeping Agent daily heartbeat

[Service]
Type=oneshot
WorkingDirectory=/home/bball_1181/thepolka.cloud/directory-app
ExecStart=/home/bball_1181/thepolka.cloud/.venv/bin/python3 housekeeping.py
